# conky/autostart/checkupdts.desktop

[Desktop Entry]
Name=data
GenericName=data
Exec=/home/sysadmin/.config/conky/scripts/checkupdts.sh
Type=Application
